# parakeet-dbus.service
[Unit]
Description=Parakeet Speech-to-Text D-Bus Service
After=sound.target

[Service]
Type=simple
ExecStartPre=/usr/bin/systemctl --user import-environment WAYLAND_DISPLAY XDG_RUNTIME_DIR
ExecStart=/usr/bin/uv run python3 -u %h/AUR/asr-hotkey/parakeet_dbus.py
WorkingDirectory=%h/AUR/asr-hotkey
EnvironmentFile=%h/.config/parakeet/parakeet.env
Environment="YDOTOOL_SOCKET=/run/ydotool/socket"
Environment="DISPLAY=:1"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
